# etc/service/cicindela_batch_track2/run
#!/bin/sh

/home/cicindela/bin/batch.pl --track=2

sleep 5
