# etc/service/cicindela_flush_buffers/run
#!/bin/sh

/home/cicindela/bin/flush_buffers.pl

sleep 5
