# etc/service/cicindela_batch_track1/run
#!/bin/sh

/home/cicindela/bin/batch.pl --track=1

sleep 5
